# celery-worker.service
[Unit]
Description=Celery Worker Service
After=network.target redis.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/path/to/braincountBackend
Environment="PATH=/path/to/braincountBackend/venv/bin"
ExecStart=/path/to/braincountBackend/venv/bin/python manage.py celery worker -l info
Restart=always
RestartSec=10

[Install]
 